# docker/supervisor/queue-worker.conf
; Mantiene activo el worker de colas (carga masiva e IA en background, RF-11).
[program:convalidaciones-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work redis --sleep=3 --tries=3 --timeout=600
autostart=true
autorestart=true
stopwaitsecs=620
numprocs=2
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/convalidaciones-worker.log
